# etc/systemd/system/mysql.service.d/override.conf
[Service]
ExecStartPre=/usr/local/bin/recover_galera_cluster /etc/mysql/conf.d/wsrep.cnf
