# etc/systemd/system/mysql.service.d/override.conf
[Service]
ExecStartPre=/usr/local/bin/recover_galera_cluster /usr/local/etc/galera_cluster_ip
